# deploy/systemd/media-lab-bot.service
[Unit]
Description=Telegram Media Lab Bot
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/media-lab
EnvironmentFile=/opt/media-lab/.env
ExecStart=/usr/bin/docker-compose -f /opt/media-lab/deploy/docker-compose.prod.yml up bot
ExecStop=/usr/bin/docker-compose -f /opt/media-lab/deploy/docker-compose.prod.yml stop bot
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

