# deploy/systemd/media-lab-worker.service
[Unit]
Description=Telegram Media Lab Image Worker
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/media-lab
EnvironmentFile=/opt/media-lab/.env
ExecStart=/usr/bin/docker-compose -f /opt/media-lab/deploy/docker-compose.prod.yml up worker-image
ExecStop=/usr/bin/docker-compose -f /opt/media-lab/deploy/docker-compose.prod.yml stop worker-image
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target


